# launchagents/com.agents.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agents.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>--directory</string>
    <string>/Users/huayang/code/agents/services/agents-mcp</string>
    <string>agents-mcp</string>
    <string>--daemon</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>8765</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/huayang/code/agents</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENTS_CONFIG_PATH</key>
    <string>/Users/huayang/code/agents/agents.yaml</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/huayang</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/huayang/code/agents/.daemon.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/huayang/code/agents/.daemon.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
