# launchagents/com.agents.admin.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agents.admin.supervisor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/huayang/code/agents/admin-supervisor.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/huayang/code/agents/.admin-supervisor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/huayang/code/agents/.admin-supervisor.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/huayang/code/agents</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/huayang</string>
    </dict>
</dict>
</plist>
