# launchagents/com.agents.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agents.telegram-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <!-- Source .env via export so the bot picks up TELEGRAM_BOT_TOKEN /
         TELEGRAM_HUMAN_CHAT_ID without baking secrets into the plist. -->
    <string>set -a; source /Users/huayang/code/agents/.env; set +a; exec /opt/homebrew/bin/uv run --directory /Users/huayang/code/agents/services/telegram-bot python bot.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/huayang/code/agents</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/huayang</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/huayang/code/agents/.telegram-bot.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/huayang/code/agents/.telegram-bot.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
