# launchagents/com.agents.console-frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.agents.console-frontend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/huayang/.nvm/versions/node/v22.15.0/bin/npm</string>
    <string>run</string>
    <string>dev</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/huayang/code/agents/apps/console/frontend</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/huayang/.nvm/versions/node/v22.15.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/huayang</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/huayang/code/agents/.console-frontend.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/huayang/code/agents/.console-frontend.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
